# hw1/.config/autostart/p5.desktop
[Desktop Entry]
Type=Application
Name=Processing
Exec=/usr/local/bin/processing-java --sketch=/home/pi/hw/hw1/processing_colorwalkers --present
